# RPi3-ubuntu-TXRX/etc/systemd/system/txrx-led.service
[Unit]
Description=USB Ethernet Activity and Error LEDs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/txrx_led.py
Restart=always
RestartSec=2
User=root
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
